# addon-radarr/radarr/rootfs/etc/s6-overlay/s6-rc.d/radarr/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Radarr
# Runs Radarr
# ==============================================================================

bashio::log.info 'Starting Radarr...'

exec /opt/bin/Radarr -nobrowser -data=/config
